# etc/replicate.service
[Unit]
Description=Replicate

[Service]
Restart=always
ExecStart=/usr/bin/replicate replicate

[Install]
WantedBy=multi-user.target
